# scripts/systemd/tumbilos-selfheal.timer
[Unit]
Description=TumbilOS Self-Heal Sensor (every 10 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
