# scripts/systemd/tumbilos-live-deploy.service
[Unit]
Description=TumbilOS Live Dashboard Deploy
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/tumbil/tumbil-os/scripts/deploy_live.sh
WorkingDirectory=%h/tumbil/tumbil-os
Environment=PATH=/home/yoygurt/.local/bin:/home/yoygurt/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/yoygurt/tumbil/dev-ops/libs:/home/yoygurt/.local/lib/python3.12/site-packages:/usr/lib/python3/dist-packages
EnvironmentFile=%h/.config/tge/tge-env
TimeoutStartSec=180

[Install]
WantedBy=default.target
