# scripts/systemd/tumbilos-priority-api.service
[Unit]
Description=TumbilOS Priority API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/tumbil/tumbil-os/scripts/tumbilos_priority_api.py
WorkingDirectory=%h/tumbil/tumbil-os
Environment=PATH=/home/yoygurt/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=%h/.config/tge/tge-env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
